# fuel_server
#!/bin/sh

. /etc/rc.subr

name=fuel_server
rcvar=fuel_server_enable

command="node /root/Fuel_Record/fuel-server/index.js" 

load_rc_config $name
run_rc_command "$1"